# deploy/log-mcp.service
[Unit]
Description=Log MCP Server - AI log analysis for BCL services
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=logmcp
Group=logmcp
WorkingDirectory=/home/logmcp/log-mcp
ExecStart=/usr/bin/node /home/logmcp/log-mcp/dist/index.js
Restart=always
RestartSec=5
SyslogIdentifier=log-mcp
Environment=NODE_ENV=production
Environment=PORT=3850
EnvironmentFile=-/home/logmcp/log-mcp/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
# Need read access to journal
SupplementaryGroups=systemd-journal

[Install]
WantedBy=multi-user.target
